# monit.timer
[Unit]
Description=Timer for monitor check

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target