# monit.service
[Unit]
Description=Get some informations on your machine

[Service]
type=oneshot
ExecStart=/usr/bin/python3 /var/monit/monit.py